# scripts/vet-nostr-announcer.service
[Unit]
Description=VET Protocol Nostr Announcer (Clean)
After=network.target

[Service]
Type=simple
# UPDATE THESE PATHS FOR YOUR DEPLOYMENT
WorkingDirectory=/path/to/Vetprotocol
ExecStart=/usr/bin/npx tsx scripts/nostr-announcer.ts
Restart=always
RestartSec=10
Environment=NODE_ENV=production
EnvironmentFile=/path/to/Vetprotocol/.env.local

[Install]
WantedBy=multi-user.target
